# etc/systemd/system/paru-sync.service
[Unit]
Description=Paru Sync Script 

[Service]
Type=simple
ExecStart=/usr/bin/paru -Sy

[Install]
WantedBy=default.target
